# job/praca/init/supervisor.conf
; This file is managed by ansible, don't make changes here - they will be overwritten.

[program:praca-autosearch]
command=/usr/bin/php /home/gutsout/h/praca-portal/backend/cli/autosearch.php
process_name=%(process_num)02d
numprocs=2
numprocs_start=1
autorestart=true
stopsignal=TERM
killasgroup=true
stdout_logfile=/var/log/supervisor/praca-autosearch.log
stderr_logfile=/var/log/supervisor/praca-autosearch.log
stdout_events_enabled=false
loglevel=warn
user=gutsout
environment=APPLICATION_ENV=development

[program:praca-autosearch-counter]
command=/usr/bin/php /home/gutsout/h/praca-portal/backend/cli/autosearch-counter.php
process_name=%(process_num)02d
numprocs=4
numprocs_start=1
autorestart=true
stopsignal=TERM
killasgroup=true
stdout_logfile=/var/log/supervisor/praca-autosearch-counter.log
stderr_logfile=/var/log/supervisor/praca-autosearch-counter.log
stdout_events_enabled=false
loglevel=warn
user=gutsout
environment=APPLICATION_ENV=development

[program:praca-search-counter]
command=/usr/bin/php /home/gutsout/h/praca-portal/backend/cli/search-counter.php
process_name=%(program_name)s
autorestart=true
stopsignal=TERM
killasgroup=true
stdout_logfile=/var/log/supervisor/praca-search-counter.log
stderr_logfile=/var/log/supervisor/praca-search-counter.log
stdout_events_enabled=false
loglevel=warn
user=gutsout
environment=APPLICATION_ENV=development

[program:praca-mailer]
command=/usr/bin/php /home/gutsout/h/praca-portal/backend/cli/mailer.php
process_name=%(program_name)s
autorestart=true
stopsignal=TERM
killasgroup=true
stdout_logfile=/var/log/supervisor/praca-mailer.log
stderr_logfile=/var/log/supervisor/praca-mailer.log
stdout_events_enabled=false
loglevel=warn
user=gutsout
environment=APPLICATION_ENV=development

[program:praca-autosearch-mailer]
command=/usr/bin/php /home/gutsout/h/praca-portal/backend/cli/mailer-autosearch.php
process_name=%(program_name)s
autorestart=true
stopsignal=TERM
killasgroup=true
stdout_logfile=/var/log/supervisor/praca-mailer-autosearch.log
stderr_logfile=/var/log/supervisor/praca-mailer-autosearch.log
stdout_events_enabled=false
loglevel=warn
user=gutsout
environment=APPLICATION_ENV=development

[program:praca-sms]
command=/usr/bin/php /home/gutsout/h/praca-portal/backend/cli/sms.php
process_name=%(program_name)s
autorestart=true
stopsignal=TERM
killasgroup=true
stdout_logfile=/var/log/supervisor/praca-sms.log
stderr_logfile=/var/log/supervisor/praca-sms.log
stdout_events_enabled=false
loglevel=warn
user=gutsout
environment=APPLICATION_ENV=development

[program:praca-push-deactivate-consumer]
command=/usr/bin/php /home/gutsout/h/praca-portal/backend/cli/push-deactive-consumer.php
process_name=%(program_name)s
autorestart=true
stopsignal=TERM
killasgroup=true
stdout_logfile=/var/log/supervisor/praca-push-deactive-consumer.log
stderr_logfile=/var/log/supervisor/praca-push-deactive-consumer.log
stdout_events_enabled=false
loglevel=warn
user=gutsout
environment=APPLICATION_ENV=development

[program:praca-push-list-consumer]
command=/usr/bin/php /home/gutsout/h/praca-portal/backend/cli/push-list-consumer.php
process_name=%(program_name)s
autorestart=true
stopsignal=TERM
killasgroup=true
stdout_logfile=/var/log/supervisor/praca-push-list-consumer.log
stderr_logfile=/var/log/supervisor/praca-push-list-consumer.log
stdout_events_enabled=false
loglevel=warn
user=gutsout
environment=APPLICATION_ENV=development

[program:praca-push-send-consumer]
command=/usr/bin/php /home/gutsout/h/praca-portal/backend/cli/push-send-consumer.php
process_name=%(program_name)s
autorestart=true
stopsignal=TERM
killasgroup=true
stdout_logfile=/var/log/supervisor/praca-push-send-consumer.log
stderr_logfile=/var/log/supervisor/praca-push-send-consumer.log
stdout_events_enabled=false
loglevel=warn
user=gutsout
environment=APPLICATION_ENV=development

[program:praca-site-list-consumer]
command=/usr/bin/php /home/gutsout/h/praca-portal/backend/cli/site-list-consumer.php
process_name=%(program_name)s
autorestart=true
stopsignal=TERM
killasgroup=true
stdout_logfile=/var/log/supervisor/praca-site-list-consumer.log
stderr_logfile=/var/log/supervisor/praca-site-list-consumer.log
stdout_events_enabled=false
loglevel=warn
user=gutsout
environment=APPLICATION_ENV=development

[program:praca-notify-router]
command=/usr/bin/php /home/gutsout/h/praca-portal/daemon/notify-router.php
process_name=%(program_name)s
autorestart=true
stopsignal=TERM
killasgroup=true
stdout_logfile=/var/log/supervisor/praca-notify-router.log
stderr_logfile=/var/log/supervisor/praca-notify-router.log
stdout_events_enabled=false
loglevel=warn
user=gutsout
environment=APPLICATION_ENV=development

[program:praca-notify-websocket]
command=/usr/bin/php /home/gutsout/h/praca-portal/daemon/notify-websocket.php
process_name=%(program_name)s
autorestart=true
stopsignal=TERM
killasgroup=true
stdout_logfile=/var/log/supervisor/praca-notify-websocket.log
stderr_logfile=/var/log/supervisor/praca-notify-websocket.log
stdout_events_enabled=false
loglevel=warn
user=gutsout
environment=APPLICATION_ENV=development

[program:praca-notify-announcer]
command=/usr/bin/php /home/gutsout/h/praca-portal/daemon/notify-announcer.php
process_name=%(program_name)s
autorestart=true
stopsignal=TERM
killasgroup=true
stdout_logfile=/var/log/supervisor/praca-notify-announcer.log
stderr_logfile=/var/log/supervisor/praca-notify-announcer.log
stdout_events_enabled=false
loglevel=warn
user=gutsout
environment=APPLICATION_ENV=development